# config/api.ini
[program:api]
directory=/app/api
command=/opt/python3.10/bin/bin/gunicorn -c prod.py cmd:app
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/api.log