# config/worker.ini
[program:worker]
directory=/opt/app/api
command=/opt/python3.10/bin/celery -A cmd.celery worker -l info -P eventlet
user=root
autostart=true
autorestart=true
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile=/var/log/supervisor/worker.log